# quadlets/home-server/botcodes/filebrowser.container
[Unit]
Description=Filebrowser Quadlet
Requires=authelia.service
After=authelia.service

[Container]
ContainerName=filebrowser
Image=docker.io/filebrowser/filebrowser
Network=botcodes.network

Volume=%h/.local/containers/filebrowser/serve:/srv
Volume=%h/.local/containers/filebrowser/database:/database
Volume=%h/.local/containers/filebrowser/config:/config:ro

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
