# quadlets/home-server/headscale.container
[Unit]
Description=Headscale Quadlet

[Container]
ContainerName=headscale
Image=docker.io/headscale/headscale:latest

Network=botcodes.network

Volume=%h/.local/containers/headscale/config:/etc/headscale:ro
Volume=%h/.local/containers/headscale/lib:/var/lib/headscale

Exec=serve

[Service]
Restart=on-failure

[Install]
WantedBy=default.target

